# backend/crontab
# Cron job to revalidate all posts where view counter changed
* * * * * curl -X POST http://localhost/blog/revalidate_views
